# services/pentryansvar_data_server.service
[Unit]
Description=Runs a server providing pentryansvar data for the school Stockholm Science & Innovation School

[Service]
ExecStart=/srv/services/start_pentryansvar_server.sh
